# supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command = gunicorn flaskr:app -b localhost:8000