# supervisord.conf
[supervisord]
nodaemon=false

[program:flaskr]
command = gunicorn flaskr:app -b localhost:8000