# supervisord.conf
[inet_http_server]
port = 127.0.0.1:9001

[supervisord]
nodaemon=false

[program:flaskr]
command = gunicorn flaskr:app -b localhost:8000