# supervisord.conf
[program:web]
command = gunicorn flaskr:app -b localhost:8000